# crontab
*/10 * * * * sh /home/joleal-b/sleep.sh; sh /home/joleal-b/monitoring.sh
